# infra/modules/ignition/base/files/containers.target
[Unit]
Description=All container services
Wants=network-online.target local-fs.target
After=network-online.target local-fs.target

[Install]
WantedBy=multi-user.target
