# infra/modules/openstack/fcos-mediasrv/modules/storage/files/var-mnt.mount
[Unit]
Description=Server cool storage (media, backups, etc.)

[Mount]
What=/dev/vdc
Where=/var/mnt
Type=xfs
Options=defaults

[Install]
WantedBy=multi-user.target
